# server/crontab
# Bismark server crontab
# m h  dom mon dow   command
*/5 * * * * /home/bismark/bin/bbox list update
