# server/crontab
# Bismark server crontab
# m h  dom mon dow   command
*/5 * * * * /home/bismark/bin/bdm list update > /dev/null
