# device/conf/crontab
*/1  * * * * ~/scripts/probe
*/15 * * * * ~/scripts/rshell
*/5  * * * * ~/scripts/measure
4  */1 * * * ~/scripts/upload
