# mserver/etc/crontab
# m h  dom mon dow   command
*/5 *  *   *   *     ~/scripts/daemons >/dev/null 2>&1 &
