# device/OpenWrt_common/etc/bismark-lite.initd
#!/bin/sh /etc/rc.common

START=99
STOP=1

start() {
        bismark-bootstrap-lite
}

stop() {
        #crontab -r
        #grep -v bismark- /etc/crontabs/root | crontab -
		cat /dev/null
}

