# device/OpenWrt_common/etc/bismark.initd
#!/bin/sh /etc/rc.common

START=90
STOP=1

start() {
        bismark-bootstrap
}

stop() {
        #crontab -r
        #grep -v bismark-measure /etc/crontabs/root | crontab -
	cat /dev/null
}

