# device/OpenWrt_common/etc/crontab
SHELL=/bin/bash
BISMARK_DIR="/tmp/bismark"

*/1  * * * * $BISMARK_DIR/scripts/probe
*/5  * * * * $BISMARK_DIR/scripts/rshell
#*/5  * * * * ~/scripts/measure
