# device/OpenWrt_common/etc/crontab
*/1  * * * * /usr/bin/bismark-probe
*/5  * * * * /usr/bin/bismark-rshell
5 */12 * * * /usr/bin/bismark-action mgmtconfupdate
