# device/OpenWrt_i586/conf/crontab
*/1  * * * * ~/scripts/probe
*/5  * * * * ~/scripts/rshell
*/5  * * * * ~/scripts/measure
4  */1 * * * ~/scripts/upload
*/10 * * * * ~/scripts/caching
